# remaster/chroot/r8723bs-bluetooth/hardware/bluetooth/scripts/bluetooth_r8723bs.rules
SUBSYSTEM=="bluetooth" RUN+="/root/hardware/bluetooth/scripts/bluetooth_udev_rules.sh"
